# target/linux/mini2440/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
tts/0::askfirst:/bin/ash --login
ttyS0::askfirst:/bin/ash --login

s3c2410_serial0::askfirst:/bin/ash --login
ttySAC0::askfirst:/bin/ash --login

